# init/airflow-web.conf
start on starting airflow
stop on stopping airflow
respawn

env AIRFLOW_HOME=...
env SENTRY_DSN=...

setuid airflow

exec airflow webserver
